# bridged.conf
[program:bridged]
command=/app/venv/bin/python manage.py runbridged
directory=/app
stopsignal=INT
stdout_logfile=/tmp/bridge.stdout.log
stderr_logfile=/tmp/bridge.stderr.log
