# celery.conf
[program:celery]
command=/home/n/vnojsite/bin/celery -A dmoj_celery worker --concurrency=1
directory=/app
stdout_logfile=/tmp/celery.stdout.log
stderr_logfile=/tmp/celery.stderr.log
