# celery.conf
[program:celery]
command=/usr/local/bin/celery -A dmoj_celery worker
directory=/app
stdout_logfile=/tmp/celery.stdout.log
stderr_logfile=/tmp/celery.stderr.log
